# files/services/sylkserver/run
#!/bin/sh

/usr/bin/sylk-server --no-fork
