# files/services/sylk/run
#!/bin/sh

cd /opt/sylk-webrtc
node_modules/.bin/gulp serve
